# cron/crontab
*/3 * * * * sh /apt/bin/download_mask_dataset.sh && python /apt/util.py
