# cron/crontab
*/3 * * * * /apt/bin/update_mask_info.sh
